# com.cyberprinter.daemon.plist

<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cyberprinter.daemon</string>
    <key>Program</key>
    <string>/Users/mars/OpenClaw_Arch_Ultra/venv/bin/python3</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mars/OpenClaw_Arch_Ultra/cyber_printer_ultimate.py</string>
        <string>--daemon</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mars/OpenClaw_Arch_Ultra</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/mars/OpenClaw_Arch_Ultra/logs/daemon.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/mars/OpenClaw_Arch_Ultra/logs/daemon.err</string>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
